# autostart/org.kde.yakuake.desktop
[Desktop Entry]
Type=Application
Name=Yakuake
Name[zh_CN]=Yakuake 下拉式终端
Comment=Start Yakuake on KDE session start
Comment[zh_CN]=在 KDE 会话启动时运行 Yakuake
Exec=yakuake
TryExec=yakuake
Icon=yakuake
Terminal=false
StartupNotify=false
DBusActivatable=true

# 仅在 KDE/Plasma 中自启动
OnlyShowIn=KDE;

# KDE 自启动调度（按照 Plasma 启动阶段与面板加载完成之后再起）
# 0: Base Desktop, 1: Desktop Services, 2: Applications
X-KDE-autostart-phase=2
X-KDE-autostart-after=panel
X-KDE-StartupNotify=false
X-KDE-SubstituteUID=false

